# 85-fpcskbd.rules
SUBSYSTEM=="hidraw", ATTRS{product}=="Computer Cool School", RUN+="/usr/local/bin/fpkbd-input /dev/$name"

