# 85-fpcskbd.rules
SUBSYSTEM=="hidraw", ATTRS{product}=="Computer Cool School", SYMLINK+="fpcskbd", RUN+="/usr/local/bin/fpkbd-input /dev/%k"
